# mic_recorder.service
[Unit]
Description=Microphone Audio Recording Service
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /usr/local/bin/mic_recorder.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mic_recorder

[Install]
WantedBy=multi-user.target 